# w.conf
[program:w]
command=python3 /home/curator2/w/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/w/test.err.log
stdout_logfile=/home/curator2/w/test.out.log